# scripts/com.tunelease.branch-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tunelease.branch-cleanup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kobayashiisaoryou/clawd/tune_lease_55/scripts/cleanup_old_branches.sh</string>
    </array>
    <!-- 毎週月曜 AM 5:00 に実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kobayashiisaoryou/Library/Logs/tunelease/branch_cleanup_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kobayashiisaoryou/Library/Logs/tunelease/branch_cleanup_stderr.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
