# launchd/com.tunelease.obsidian-reindex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tunelease.obsidian-reindex</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/kobayashiisaoryou/clawd/tune_lease_55/scripts/run_obsidian_reindex.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kobayashiisaoryou/clawd/tune_lease_55</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OBSIDIAN_VAULT_PATH</key>
    <string>/Users/kobayashiisaoryou/Library/Mobile Documents/iCloud~md~obsidian/Documents/Obsidian Vault</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>PATH</key>
    <string>/Users/kobayashiisaoryou/clawd/tune_lease_55/.venv/bin:/Users/kobayashiisaoryou/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- 毎晩 03:00 に実行（reindex完了後にChromaDB→GCS sync + Cloud Run再起動） -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kobayashiisaoryou/Library/Logs/tune_lease_55_obsidian_reindex.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kobayashiisaoryou/Library/Logs/tune_lease_55_obsidian_reindex.err.log</string>
</dict>
</plist>
