# launchd/com.tunelease.weekly-health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tunelease.weekly-health-check</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/kobayashiisaoryou/anaconda3/bin/python3</string>
    <string>/Users/kobayashiisaoryou/clawd/tune_lease_55/scripts/check_system_health.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kobayashiisaoryou/clawd/tune_lease_55</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OBSIDIAN_VAULT</key>
    <string>/Users/kobayashiisaoryou/Library/Mobile Documents/iCloud~md~obsidian/Documents/Obsidian Vault</string>
    <key>OBSIDIAN_BACKUP_ROOT</key>
    <string>/Users/kobayashiisaoryou/Library/Mobile Documents/com~apple~CloudDocs/tune_lease_55_backups/obsidian</string>
    <key>CASE_DATA_BACKUP_ROOT</key>
    <string>/Users/kobayashiisaoryou/Library/Mobile Documents/com~apple~CloudDocs/tune_lease_55_backups/case_data</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>PATH</key>
    <string>/Users/kobayashiisaoryou/anaconda3/bin:/Users/kobayashiisaoryou/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- 毎週日曜 05:00 に実行（04:00 の Obsidian バックアップ後） -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/kobayashiisaoryou/Library/Logs/tune_lease_55_weekly_health.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kobayashiisaoryou/Library/Logs/tune_lease_55_weekly_health.err.log</string>
</dict>
</plist>
